# deploy/com.vizdashboard.projects-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.vizdashboard.projects-sync</string>
  <!-- Dumb file copy: ship Claude Code transcripts to the Pi. No API key, no LLM here.
       Uses a dedicated passphraseless key so it runs headlessly from launchd.
       Destination is the `monitoring-pi` SSH alias (in ~/.ssh/config), NOT a hardcoded
       IP — the Pi's DHCP lease changes on reboot, so following the alias means a future
       IP change is a one-line ~/.ssh/config edit instead of editing this plist. -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/rsync</string>
    <string>-az</string>
    <string>--delete</string>
    <string>--include=*/</string>
    <string>--include=*.jsonl</string>
    <string>--exclude=*</string>
    <string>-e</string>
    <string>ssh -i /Users/YOUR_USERNAME/.ssh/viz_pi_sync -o IdentitiesOnly=yes -o StrictHostKeyChecking=accept-new</string>
    <string>/Users/YOUR_USERNAME/.claude/projects/</string>
    <string>monitoring-pi:claude-projects/</string>
  </array>
  <key>StartInterval</key><integer>1200</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardErrorPath</key><string>/tmp/viz-projects-sync.err</string>
  <key>StandardOutPath</key><string>/tmp/viz-projects-sync.out</string>
</dict>
</plist>
